# rootfs/etc/inittab
::sysinit:/etc/init.d/rcS
ttyO0::once:-/bin/login
